# deploy/sportplatz-buchung.service
[Unit]
Description=Sportplatz Buchungssystem (Web/Admin)
After=network.target
StartLimitBurst=3
StartLimitIntervalSec=120
OnFailure=sportplatz-crash@%n.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/git.com/Sportplatz-Buchung
ExecStart=/root/git.com/Sportplatz-Buchung/.venv/bin/uvicorn web.main:app --host 0.0.0.0 --port 1946
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
